# deploy/systemd/pidvn25006.service
[Unit]
Description=PIDVN25006 AGV Vision Runtime Supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/PIDVN25006
ExecStartPre=/usr/bin/mkdir -p /opt/PIDVN25006/outputs/runtime/supervisor
ExecStart=/opt/PIDVN25006/run_forever.sh --no-frontend
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20
StandardOutput=append:/opt/PIDVN25006/outputs/runtime/supervisor/systemd_stdout.log
StandardError=append:/opt/PIDVN25006/outputs/runtime/supervisor/systemd_stderr.log

[Install]
WantedBy=multi-user.target
